# deploy/gcp/box.service
[Unit]
Description=box API service
After=network.target chromium.service
Requires=chromium.service

[Service]
Type=simple
User=box
WorkingDirectory=/opt/box
EnvironmentFile=/etc/box/box.env
ExecStart=/opt/box/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8080 --workers 1
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
