# deploy/gcp/chromium.service
[Unit]
Description=Chromium headless (CDP) for box
After=network.target

[Service]
Type=simple
User=box
WorkingDirectory=/opt/box
EnvironmentFile=/etc/box/box.env
ExecStart=/usr/local/bin/start-chromium
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
